# systemd/.config/systemd/user/change-wallpaper.timer
[Unit]
Description=Change wallpaper everyday 

[Timer]
#Daily at 6:05AM
OnCalendar=*-*-* 06:05:00
Persistent=True

[Install]
WantedBy=timers.target
