# systemd/.config/systemd/user/on-logout.service

[Unit]
Description=Execute on logout

[Service]
Type=oneshot
ExecStart=
ExecStop=/usr/bin/bash %h/Scripts/on-logout.sh
RemainAfterExit=yes

