# systemd/.config/systemd/user/on-logout.service

[Unit]
Description=Execute on logout

[Service]
Type=oneshot
ExecStart=
ExecStop=/usr/bin/task +ACTIVE stop
RemainAfterExit=yes

